# infra/systemd/nomaddev-orchestrator.service
[Unit]
Description=NomadDev orchestrator (Tailscale-fronted WebSocket relay)
Documentation=https://github.com/MattCheramie/NomadDev
# Wait for the network stack so /healthz binds successfully on boot.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nomaddev
Group=nomaddev
EnvironmentFile=/etc/nomaddev/env
ExecStart=/usr/local/bin/orchestrator -listen :8080
# Restart=always (not on-failure): the orchestrator exits cleanly (status 0)
# when the web settings editor requests a config-change restart, and that
# must bounce the service too. `systemctl stop` still stops it — systemd
# tells an operator stop apart from a process exit.
Restart=always
RestartSec=2s

# Hardening. The orchestrator does not need root, special caps, or write
# access outside /var/lib/nomaddev (sessions.db, history.db, workspace).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
ReadWritePaths=/var/lib/nomaddev
# Stdout/stderr → journald (the binary already writes JSON via slog).
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
